# openwrt-feed/net/nordvpn-wireguard/files/etc/uci-defaults/99-nordvpn-cron
#!/bin/sh

/etc/init.d/nordvpn enable >/dev/null 2>&1
/etc/init.d/nordvpn reload >/dev/null 2>&1

exit 0
